# gauge/gauge.service
[Unit]
Description=Gauge
After=multi-user.target

[Service]
User=linkage
ExecStart=/usr/bin/gauge
Restart=always
RestartSec=3s

[Install]
WantedBy=multi-user.target